# system/g15-daemon.service
[Unit]
Description=Dell G15 Control Center Daemon
After=multi-user.target
ConditionPathExists=/proc/acpi/call

[Service]
Type=simple
User=root
ExecStart=/opt/g15-controller/venv/bin/python /opt/g15-controller/src/g15_daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target